# .config/containers/systemd/vaultwarden.container
[Container]
ContainerName=vaultwarden
Image=docker.io/vaultwarden/server:latest
AutoUpdate=registry
#Environment=ROCKET_PORT=8080
Volume=vaultwarden-data:/data:Z
Network=oci.network

[Service]
Restart=always

[Install]
WantedBy=default.target
